# crontab
* * * * * /bin/sh /app/check_log.sh >> /var/log/cron.log 2>&1
